# stage-experiment/src/main/resources/ansible/background-traffic-driver.service
[Unit]
Description=MAP Background Traffic Driver
Requires=network.target
After=network.target

[Service]
User=map
WorkingDirectory=/var/lib/map/background-traffic
ExecStart=/usr/bin/java -Xmx1G -jar /var/lib/map/background-traffic/background-traffic-driver.jar
SuccessExitStatus=0
Restart=on-failure
RestartSec=5s
OOMScoreAdjust=-1000

[Install]
WantedBy=default.target
