# stage-experiment/src/main/resources/ansible/map-agent.service
[Unit]
Description=MAP Agent
Requires=network.target
After=network.target

[Service]
User=map
WorkingDirectory=/var/lib/map/agent
ExecStart=/usr/bin/java -Xmx1G -jar /var/lib/map/agent/map-agent.jar
SuccessExitStatus=0
Restart=on-failure
RestartSec=5s
OOMScoreAdjust=-1000

[Install]
WantedBy=default.target
